# combustion/files/poweroff.timer
[Unit]
Description=Poweroff daily

[Timer]
OnCalendar=*-*-* 00:10:*
Unit=poweroff.service

[Install]
WantedBy=timers.target
