# combustion/files/poweroff.service
[Unit]
Description=Poweroff job for timer

[Service]
Type=oneshot
ExecStart=systemctl poweroff
