# vpn-client/etc/systemd/system/nftables.service
[Unit]
Description=Initialises nftab rules

[Service]
Type=oneshot
ExecStart=/etc/nftables.conf

[Install]
WantedBy=default.target
