# combustion/apps/system/files/home-cloud-backup.service
[Unit]
Description="Home Cloud Backup"

[Service]
ExecStart=/cloud/system/backup.sh
